# deploy/sessionrunner-stop.service
[Unit]
Description=Arrete le container SessionRunner

[Service]
Type=oneshot
ExecStart=/usr/bin/docker stop -t 15 sessionrunner
